# deploy/systemd/atmos@.service
[Unit]
Description=AirGradient Sensor Collector - %i
After=network.target influxdb.service
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=5
User=USER_PLACEHOLDER
WorkingDirectory=DIR_PLACEHOLDER

# 1. Load global defaults (Influx tokens, etc.)
EnvironmentFile=-DIR_PLACEHOLDER/.env

# 2. Load instance-specific config (SENSOR_IP or SENSOR_SERIAL)
EnvironmentFile=-/etc/atmos/%i.env

# 3. Identity via Flags
# InfluxDB/Global credentials come from .env, Identity comes from flags here.
# Note: ${VAR} syntax allows using variables loaded from EnvironmentFiles.
ExecStart=/usr/local/bin/atmos collect \
  --location %i \
  --ip ${SENSOR_IP} \
  --serial ${SENSOR_SERIAL} \
  --interval 1m

[Install]
WantedBy=multi-user.target
